# maintainer/testdirs/testdirl2/autostart/tdg1-1.desktop
[Desktop Entry]
Name=tdg1.1
Exec=touch tdg1.1-fake
