# maintainer/testdirs/testdirl2/autostart/tdg1-2.desktop
[Desktop Entry]
Name=tdg1.2
Exec=touch tdg1.2
