# maintainer/testdirs/testdirl2/autostart/tdl2-1.desktop
[Desktop Entry]
Name=tdl2.1
Exec=touch tdl2.1
