# maintainer/testdirs/testdirg2/autostart/tdg2-2.desktop
[Desktop Entry]
Name=tdg2.2
Exec=touch tdg2.2
