# maintainer/testdirs/testdirg2/autostart/tdg2-1.desktop
[Desktop Entry]
Name=tdg2.1
Exec=touch tdg2.1
