# maintainer/testdirs/testdirg1/autostart/tdl1-1.desktop
[Desktop Entry]
Name=tdl1.1
Exec=touch tdl1.1
