# maintainer/testdirs/testdirg3/autostart/tdg3-2.desktop
[Desktop Entry]
Name=tdg3.2
Exec=touch tdg3.2
