# maintainer/testdirs/testdirg3/autostart/tdg3-1.desktop
[Desktop Entry]
Name=tdg3.1
Exec=touch tdg3.1
